# usr-local/lib/systemd/system/clash-renew.service
[Unit]
Description=renew clash config via subscribe url

[Service]
Type=oneshot
User=root
ExecStart=clash-renew
PrivateTmp=true
